# packaging/systemd/e-clawhisper.service
[Unit]
Description=e-clawhisper voice daemon
Documentation=https://github.com/eager-dev/e-clawhisper
After=sound.target pipewire.service pipewire-pulse.service
Wants=pipewire.service

[Service]
Type=simple

EnvironmentFile=-%h/.config/e-clawhisper/.env
EnvironmentFile=-/etc/e-clawhisper/.env
Environment=MODELS_DIR=/opt/e-clawhisper/models
Environment=DATA_DIR=/opt/e-clawhisper/data
Environment=SOCKET_PATH=%t/e-clawhisper/eclaw.sock
Environment=PYTHONPATH=/opt/e-clawhisper/src

ExecStartPre=/bin/mkdir -p %t/e-clawhisper
ExecStart=/opt/e-clawhisper/.venv/bin/python -m e_clawhisper.main session start
ExecStop=/usr/bin/eclaw session stop

Restart=on-failure
RestartSec=5
TimeoutStopSec=10

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
PrivateTmp=yes
ReadWritePaths=%t/e-clawhisper

[Install]
WantedBy=default.target
